# deploy/ansible/roles/ingraind/templates/etc/systemd/system/ingraind.service.j2
[Unit]
Description=InGRAINd data collection agent
Requires=network.service

[Service]
ExecStart=/usr/local/bin/ingraind /etc/ingraind/ingraind.toml